# rootfs/etc/udev/rules.d/85-usb-flash.rules
# do not edit this file, it will be overwritten on update

# run a command on remove events
SUBSYSTEM=="block", ACTION=="remove", RUN+="/usr/loader/reboot.sh"
SUBSYSTEM=="block", ACTION=="add", RUN+="/usr/loader/loader_32.sh"

