# rootfs/etc/udev/rules.d/97-usb-keyboard.rules
# do not edit this file, it will be overwritten on update

# run a command on remove events
BUS=="usb", SUBSYSTEM=="input", ACTION=="add", RUN+="/etc/X11/Xinit.d/12keymap"
