# rootfs_new/usr/lib/opkg/info/avahi-daemon.list
/usr/share/avahi/introspection/DomainBrowser.introspect
/usr/share/avahi/introspection/RecordBrowser.introspect
/usr/share/avahi/introspection/ServiceResolver.introspect
/usr/share/avahi/introspection/ServiceBrowser.introspect
/usr/sbin/avahi-daemon
/usr/share/avahi/introspection/AddressResolver.introspect
/etc/dbus-1/system.d/avahi-dbus.conf
/usr/share/avahi/avahi-service.dtd
/usr/share/avahi/introspection/EntryGroup.introspect
/etc/init.d/avahi-daemon
/etc/avahi/avahi-daemon.conf
/usr/share/avahi/introspection/Server.introspect
/usr/share/avahi/introspection/ServiceTypeBrowser.introspect
/etc/avahi/hosts
/usr/share/avahi/introspection/HostNameResolver.introspect
/usr/share/avahi/service-types
/etc/avahi/services/ssh.service
/etc/avahi/services/sftp-ssh.service
