# rootfs_new/usr/lib/opkg/info/ipaq-sleep.list
/etc/ipaq-sleep.conf
/etc/X11/Xsession.d/69ipaq-sleep
/etc/init.d/ipaq-sleep
/usr/bin/ipaq-sleep
