# rootfs_new/usr/lib/opkg/info/netplug.prerm
#!/bin/sh
set -e
# Automatically added by dh_installinit
if [ -x "/etc/init.d/netplug" ]; then
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d netplug stop || exit $?
	else
		/etc/init.d/netplug stop || exit $?
	fi
fi
# End automatically added section
