# rootfs_new/usr/lib/opkg/info/netplug.conffiles
/etc/init.d/netplug
/etc/netplug/netplugd.conf
/etc/netplug/netplug
