# rootfs_new/usr/lib/opkg/info/usb-gadget-mode.list
/etc/init.d/usb-gadget
/usr/bin/usb-gadget
/etc/default/usb-gadget
