# rootfs_new/usr/lib/opkg/info/usb-gadget-mode.postinst
	test -n "$D" && opt="-r $D"
	update-rc.d $opt usb-gadget defaults 50
