# rootfs_new/etc/udev/rules.d/99-usb-mouse.rules
# do not edit this file, it will be overwritten on update

# run a command on remove events
ACTION=="add", SUBSYSTEM=="input", RUN+="/usr/mouse_ptr/show_ptr.sh"
ACTION=="remove", SUBSYSTEM=="input", RUN+="/usr/mouse_ptr/hide_ptr.sh"

