# rootfs/usr/lib/opkg/info/busybox-syslog.postrm
#!/bin/sh
if test "x$D" != "x"; then
	OPT="-r $D"
else
	OPT=""
fi
update-rc.d $OPT syslog remove
