# rootfs/usr/lib/opkg/info/busybox.list
/etc/busybox.links
/etc/udhcpc.d/50default
/bin/sh
/bin/busybox
/usr/share/udhcpc/default.script
/etc/init.d/hwclock.sh
