# rootfs/usr/lib/opkg/info/netplug.list
/usr/share/doc/netplug/TODO
/usr/share/man/man8/netplugd.8.gz
/usr/share/doc/netplug/NEWS.gz
/usr/share/doc/netplug/README
/etc/netplug/netplug
/etc/netplug/netplugd.conf
/usr/share/doc/netplug/copyright
/usr/share/doc/netplug/changelog.gz
/etc/init.d/netplug
/usr/share/doc/netplug/changelog.Debian.gz
/sbin/netplugd
