# rootfs/usr/lib/opkg/info/openssh-sshd.list
/etc/ssh/moduli
/usr/sbin/sshd
/etc/ssh/sshd_config
/etc/init.d/sshd
