# rootfs/usr/lib/opkg/info/openssh-sshd.postrm
if test "x$D" != "x"; then
	exit 1
else
	/etc/init.d/sshd stop
	deluser sshd
	delgroup sshd
	update-rc.d -f sshd remove
fi
