# rootfs/usr/lib/opkg/info/procps.list
/usr/bin/w
/usr/bin/watch
/bin/kill.procps
/usr/bin/pmap
/usr/bin/pwdx
/usr/bin/snice
/usr/bin/pgrep.procps
/usr/bin/top.procps
/usr/bin/vmstat
/etc/sysctl.conf
/usr/bin/slabtop
/etc/init.d/procps.sh
/usr/bin/free.procps
/usr/bin/skill
/usr/bin/tload
/bin/ps.procps
/sbin/sysctl.procps
/usr/bin/uptime.procps
/lib/libproc-3.2.8.so
/usr/bin/pkill.procps
