# rootfs/usr/lib/opkg/info/sysvinit.list
/sbin/halt.sysvinit
/etc/default/rcS
/sbin/poweroff.sysvinit
/etc/init.d/rc
/bin/mountpoint
/etc/init.d/rcS
/sbin/reboot.sysvinit
/sbin/telinit
/sbin/init.sysvinit
/sbin/killall5
