# rootfs/usr/lib/opkg/info/udev.list
/lib/udev/cdrom_id
/etc/udev/rules.d/64-md-raid.rules
/etc/udev/rules.d/60-persistent-v4l.rules
/etc/udev/rules.d/60-persistent-storage-tape.rules
/etc/udev/rules.d/local.rules
/etc/udev/rules.d/60-cdrom_id.rules
/lib/udev/fstab_import
/etc/udev/rules.d/40-ia64.rules
/lib/udev/vol_id
/etc/udev/saved.atags
/lib/libudev.so.0.2.0
/etc/udev/links.conf
/lib/udev/edd_id
/lib/udev/ata_id
/etc/udev/rules.d/40-pilot-links.rules
/etc/init.d/udev
/etc/udev/saved.uname
/etc/udev/rules.d/79-fstab_import.rules
/lib/udev/path_id
/etc/udev/rules.d/udev.rules
/etc/udev/rules.d/40-zaptel.rules
/lib/udev/rules.d
/etc/udev/rules.d/40-s390.rules
/lib/udev/create_floppy_devices
/etc/udev/rules.d/permissions.rules
/etc/udev/rules.d/75-persistent-net-generator.rules
/etc/udev/rules.d/61-persistent-storage-edd.rules
/etc/udev/rules.d/60-persistent-serial.rules
/lib/udev/write_net_rules
/etc/udev/rules.d/40-infiniband.rules
/etc/udev/rules.d/run.rules
/lib/udev/firmware.sh
/etc/udev/rules.d/95-udev-late.rules
/etc/udev/rules.d/40-isdn.rules
/lib/udev/usb_id
/etc/udev/scripts/network.sh
/lib/udev/rule_generator.functions
/etc/udev/udev.conf
/etc/udev/rules.d/75-cd-aliases-generator.rules
/etc/udev/rules.d/50-udev-default.rules
/etc/udev/rules.d/50-firmware.rules
/lib/udev/scsi_id
/etc/udev/rules.d/40-ppc.rules
/etc/scsi_id.config
/lib/udev/write_cd_rules
/sbin/udevd
/lib/udev/collect
/etc/udev/mount.blacklist
/etc/udev/scripts/mount.sh
/etc/udev/rules.d/60-persistent-input.rules
/etc/default/udev
/lib/libudev.so.0
/etc/udev/rules.d/64-device-mapper.rules
/etc/udev/saved.cmdline
/etc/udev/rules.d/40-alsa.rules
/etc/udev/rules.d/60-persistent-storage.rules
/etc/udev/rules.d/80-drivers.rules
